# va_parser/attack-technique/T1546.002.json
{
    "id": "T1546.002",
    "name": "Screensaver",
    "x_mitre_is_subtechnique": true,
    "revoked": false,
    "x_mitre_deprecated": false,
    "description": "Adversaries may establish persistence by executing malicious content triggered by user inactivity. Screensavers are programs that execute after a configurable time of user inactivity and consist of Portable Executable (PE) files with a .scr file extension.(Citation: Wikipedia Screensaver) The Windows screensaver application scrnsave.scr is located in <code>C:\\Windows\\System32\\</code>, and <code>C:\\Windows\\sysWOW64\\</code>  on 64-bit Windows systems, along with screensavers included with base Windows installations.\n\nThe following screensaver settings are stored in the Registry (<code>HKCU\\Control Panel\\Desktop\\</code>) and could be manipulated to achieve persistence:\n\n* <code>SCRNSAVE.exe</code> - set to malicious PE path\n* <code>ScreenSaveActive</code> - set to '1' to enable the screensaver\n* <code>ScreenSaverIsSecure</code> - set to '0' to not require a password to unlock\n* <code>ScreenSaveTimeout</code> - sets user inactivity timeout before screensaver is executed\n\nAdversaries can use screensaver settings to maintain persistence by setting the screensaver to run malware after a certain timeframe of user inactivity.(Citation: ESET Gazer Aug 2017)",
    "x_mitre_platforms": [
        "Windows"
    ],
    "x_mitre_domains": [
        "enterprise-attack"
    ],
    "x_mitre_defense_bypassed": [],
    "external_references": [
        {
            "source_name": "mitre-attack",
            "url": "https://attack.mitre.org/techniques/T1546/002",
            "external_id": "T1546.002"
        },
        {
            "source_name": "ESET Gazer Aug 2017",
            "description": "ESET. (2017, August). Gazing at Gazer: Turla\u2019s new second stage backdoor. Retrieved September 14, 2017.",
            "url": "https://www.welivesecurity.com/wp-content/uploads/2017/08/eset-gazer.pdf"
        },
        {
            "source_name": "Wikipedia Screensaver",
            "description": "Wikipedia. (2017, November 22). Screensaver. Retrieved December 5, 2017.",
            "url": "https://en.wikipedia.org/wiki/Screensaver"
        }
    ],
    "kill_chain_phases": [
        {
            "kill_chain_name": "mitre-attack",
            "phase_name": "privilege-escalation"
        },
        {
            "kill_chain_name": "mitre-attack",
            "phase_name": "persistence"
        }
    ],
    "tactic": [
        "privilege-escalation",
        "persistence"
    ],
    "x_mitre_detection": "Monitor process execution and command-line parameters of .scr files. Monitor changes to screensaver configuration changes in the Registry that may not correlate with typical user behavior.\n\nTools such as Sysinternals Autoruns can be used to detect changes to the screensaver binary path in the Registry. Suspicious paths and PE files may indicate outliers among legitimate screensavers in a network and should be investigated.",
    "x_mitre_data_sources": [
        "File: File Creation",
        "Process: Process Creation",
        "Command: Command Execution",
        "Windows Registry: Windows Registry Key Modification",
        "File: File Modification"
    ],
    "x_mitre_permissions_required": [
        "User"
    ],
    "x_mitre_effective_permissions": [],
    "x_mitre_system_requirements": [],
    "x_mitre_remote_support": false,
    "relation": {
        "mitigates": [
            {
                "name": "Execution Prevention",
                "description": "Block execution of code on a system through application control, and/or script blocking.",
                "type": "course-of-action",
                "id": "M1038"
            },
            {
                "name": "Disable or Remove Feature or Program",
                "description": "Remove or deny access to unnecessary and potentially vulnerable software to prevent abuse by adversaries.",
                "type": "course-of-action",
                "id": "M1042"
            }
        ]
    },
    "mitigationRequirements": {
        "requirements": [
            {
                "description": "Ensure that Group Policy is implemented to disable screensavers if not required",
                "acid": "98580650"
            },
            {
                "description": "Ensure that .scr files are blocked from being executed from non-standard locations",
                "acid": "98580651"
            }
        ]
    },
    "detectionRequirements": {
        "requirements": [
            {
                "description": "Ensure that process execution and command-line parameters of .scr files are monitored",
                "acid": "98580646"
            },
            {
                "description": "Ensure that changes to screensaver configuration changes in the Registry are monitored",
                "acid": "98580647"
            },
            {
                "description": "Ensure that the Sysinternals Autoruns tool is deployed to detect changes to the screensaver binary path in the Registry",
                "acid": "98580648"
            },
            {
                "description": "Ensure that paths and PE files in a network should be investigated",
                "acid": "98580649"
            }
        ]
    },
    "cspcontrolIds": [],
    "externalcontrolIds": [],
    "peoplecontrolIds": [],
    "policycontrolIds": [],
    "technologycontrolIds": [
        "15010333",
        "15010334",
        "15010335",
        "15010336",
        "15020336",
        "15020337",
        "15020338",
        "15020339",
        "15030307",
        "15030308",
        "15030309",
        "15030310",
        "15040316",
        "15040317",
        "15040320",
        "15040321",
        "15180157",
        "15180158",
        "15180159",
        "15200325",
        "15200326",
        "15200327",
        "15200328",
        "16010293",
        "16010294",
        "16010295",
        "16010296",
        "16020301",
        "16020302",
        "16020303",
        "16020304",
        "16100340",
        "16100341",
        "16100342",
        "16100343",
        "16110280",
        "16110281",
        "16110282",
        "16110283",
        "16120296",
        "16120297",
        "16120298",
        "16120299",
        "16140157",
        "16140158",
        "16140159",
        "16360361",
        "16360362",
        "16360363",
        "16360364",
        "16460356",
        "16460357",
        "16460358",
        "19920276",
        "19920277",
        "19920278"
    ],
    "cveIds": []
}
