# va_parser/attack-technique/T1547.011.json
{
    "id": "T1547.011",
    "name": "Plist Modification",
    "x_mitre_is_subtechnique": true,
    "revoked": true,
    "x_mitre_deprecated": false,
    "description": "Adversaries can modify property list files (plist files) to execute their code as part of establishing persistence. Plist files are used by macOS applications to store properties and configuration settings for applications and services. Applications use information plist files, <code>Info.plist</code>, to tell the operating system how to handle the application at runtime using structured metadata in the form of keys and values. Plist files are formatted in XML and based on Apple's Core Foundation DTD and can be saved in text or binary format.(Citation: fileinfo plist file description) \n\nAdversaries can modify paths to executed binaries, add command line arguments, and insert key/pair values to plist files in auto-run locations which execute upon user logon or system startup. Through modifying plist files in these locations, adversaries can also execute a malicious dynamic library (dylib) by adding a dictionary containing the <code>DYLD_INSERT_LIBRARIES</code> key combined with a path to a malicious dylib under the <code>EnvironmentVariables</code> key in a plist file. Upon user logon, the plist is called for execution and the malicious dylib is executed within the process space. Persistence can also be achieved by modifying the <code>LSEnvironment</code> key in the application's <code>Info.plist</code> file.(Citation: wardle artofmalware volume1)",
    "x_mitre_platforms": [
        "macOS"
    ],
    "x_mitre_domains": [
        "enterprise-attack"
    ],
    "x_mitre_defense_bypassed": [],
    "external_references": [
        {
            "source_name": "mitre-attack",
            "url": "https://attack.mitre.org/techniques/T1547/011",
            "external_id": "T1547.011"
        },
        {
            "source_name": "fileinfo plist file description",
            "description": "FileInfo.com team. (2019, November 26). .PLIST File Extension. Retrieved October 12, 2021.",
            "url": "https://fileinfo.com/extension/plist"
        },
        {
            "source_name": "wardle artofmalware volume1",
            "description": "Patrick Wardle. (2020, August 5). The Art of Mac Malware Volume 0x1: Analysis. Retrieved March 19, 2021.",
            "url": "https://taomm.org/vol1/pdfs.html"
        }
    ],
    "kill_chain_phases": [
        {
            "kill_chain_name": "mitre-attack",
            "phase_name": "persistence"
        },
        {
            "kill_chain_name": "mitre-attack",
            "phase_name": "privilege-escalation"
        }
    ],
    "tactic": [
        "persistence",
        "privilege-escalation"
    ],
    "x_mitre_detection": "Monitor for common command-line editors used to modify plist files located in auto-run locations, such as <code>~/LaunchAgents</code>, <code>~/Library/Application Support/com.apple.backgroundtaskmanagementagent/backgrounditems.btm</code>, and an application's <code>Info.plist</code>. \n\nMonitor for plist file modification immediately followed by code execution from <code>~/Library/Scripts</code> and <code>~/Library/Preferences</code>. Also, monitor for significant changes to any path pointers in a modified plist.\n\nIdentify new services executed from plist modified in the previous user's session. ",
    "x_mitre_data_sources": [
        "File: File Modification",
        "Service: Service Creation",
        "Command: Command Execution",
        "Process: Process Creation"
    ],
    "x_mitre_permissions_required": [
        "User",
        "Administrator"
    ],
    "x_mitre_effective_permissions": [],
    "x_mitre_system_requirements": [],
    "x_mitre_remote_support": false,
    "relation": {
        "mitigates": []
    },
    "mitigationRequirements": {
        "requirements": [
            {
                "description": "Ensure that \"read-only\" is set for plist files to prevent them from being modified by users",
                "acid": "98580521"
            },
            {
                "description": "Ensure that users are suggested to hold the Shift key while logging to prevents apps from opening automatically",
                "acid": "98580522"
            }
        ]
    },
    "detectionRequirements": {
        "requirements": [
            {
                "description": "Ensure that file system monitoring is performed to detect modified plist files",
                "acid": "98580515"
            },
            {
                "description": "Ensure that the Knock Knock tool is deployed to detect persistence mechanisms and the specific files that are being referenced",
                "acid": "98580516"
            },
            {
                "description": "Ensure that Login items located in Apple menu/System Preferences/Users & Groups/Login items should be monitored",
                "acid": "98580517"
            },
            {
                "description": "Ensure that Login items located in Contents/Library/LoginItems within an application bundle should be monitored",
                "acid": "98580518"
            },
            {
                "description": "Ensure that process execution is monitored for abnormal process execution resulting from modified plist files",
                "acid": "98580519"
            },
            {
                "description": "Ensure that utilities used to modify plist files or that take a plist file as an argument, are monitored",
                "acid": "98580520"
            }
        ]
    },
    "cspcontrolIds": [],
    "externalcontrolIds": [],
    "peoplecontrolIds": [],
    "policycontrolIds": [],
    "technologycontrolIds": [],
    "cveIds": []
}
